# crontab
0 * * * * XAUTHORITY="/home/idegtiarenko/.Xauthority" DISPLAY=":0.0" XDG_RUNTIME_DIR="/run/user/1000" /usr/bin/notify-send -i mate-eyes-applet "Eye exercise time"
